# jobs.cron
MAILTO=""
@reboot echo "$HOME/bin/tasks.sh > /dev/null" | /usr/bin/at 'now + 2 minutes'