# jobs.cron
MAILTO=""
# Remove "stale" inhibit file
@reboot rm -f /tmp/INHIBIT
# Remove old AT jobs from the queue
@reboot atq -q = | cut -f1 | xargs atrm 1> /dev/null 2>&1
# Start the tasks immediately or delay to allow user login
@reboot $HOME/bin/modecheck.sh 1> /dev/null 2>&1
