# jobs.cron
MAILTO=""
@reboot echo "$HOME/bin/tasks.sh 1> /dev/null 2>&1" | /usr/bin/at 'now + 2 minutes'